# templates/ebs-snapshot.cron.j2
#!/usr/bin/env bash

/opt/oulib/aws/bin/ebs-snapshot.sh {{ ebs_snapshot_retention_days }}
